# systemd/user/swww-daemon.service
[Unit]
Description=swww daemon for wallpaper management
After=graphical-session.target
Requires=graphical-session.target

[Service]
ExecStart=/usr/bin/swww-daemon
Restart=on-failure
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-1

[Install]
WantedBy=graphical-session.target

